# debian/postinst
#!/bin/bash

# add openrqm user to database

# execute sql file to local database

#enable systemd service unit
systemctl enable openrqm-server.service
